# medc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

python3 /root/medc/check.py &
python3 /root/medc/meter/main.py

exit 0
